# etc/supervisor/conf.d/vncserver.conf
[program:vncserver]
command=vncserver -geometry 1280x1024 :1
user=root
stdout_logfile=/var/log/novnc.log
redirect_stderr=true
