# etc/supervisor/conf.d/vncserver.conf
[program:vncserver]
command=vncserver :1 -depth 24 -geometry 1360x768
user=root
stdout_logfile=/var/log/novnc.log
redirect_stderr=true
