# etc/supervisor/conf.d/vncserver.conf
[program:vncserver]
command=vncserver :1 -geometry 1024x768 
user=root
stdout_logfile=/var/log/novnc.log
redirect_stderr=true
