# etc/supervisor/conf.d/web.conf
[program:noVNC]
command=websockify -D --web=/usr/share/novnc/  6901 localhost:5901
user=root
stdout_logfile=/var/log/websockify.log
redirect_stderr=true
